# deploy/stocks-web.service
[Unit]
Description=stocks.playmations.com web (Kestrel)
After=network.target

[Service]
User=stocks
WorkingDirectory=/var/www/stocks/web
ExecStart=/usr/bin/dotnet /var/www/stocks/web/StocksWeb.dll
Restart=always
RestartSec=5
Environment=ASPNETCORE_URLS=http://127.0.0.1:5000
Environment=DOTNET_ENVIRONMENT=Production
Environment=DbPath=/var/lib/stocks/stocks.db
# Für On-Demand-Backtests (worker/backtest.py wird vom Web gespawnt, liest read-only)
Environment=PythonPath=/var/www/stocks/venv/bin/python
# AdminPassword etc. — nicht ins Repo, nur auf dem Server:
EnvironmentFile=-/etc/stocks/env

[Install]
WantedBy=multi-user.target
